# rc/zsh/profile.d/nodebrew.zsh
export PATH=$PATH:$HOME/.nodebrew/current/bin
ulimit -n 2560
